# src/gnome-terminal-server.service.in
[Unit]
Description=GNOME Terminal Server

[Service]
KillMode=process
Type=dbus
BusName=org.gnome.Terminal
ExecStart=@libexecdir@/gnome-terminal-server
